# systemd-services/conversation-service.service
[Unit]
Description=Harena Conversation Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/harena
Environment="PATH=/home/ubuntu/harena/env/bin"
EnvironmentFile=/home/ubuntu/harena/.env
ExecStart=/home/ubuntu/harena/env/bin/uvicorn conversation_service.main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
